# orchestrator/gha-orch-trigger.service
[Unit]
Description=Reconcile runner slots on LAN trigger (socket-activated)
Requires=gha-orch-trigger.socket

[Service]
Type=oneshot
ExecStart=/opt/gha-local-orchestrator/gha-orch-trigger.sh
