# orchestrator/gha-local-orchestrator.timer
[Unit]
Description=Run GitHub Actions local runner orchestrator periodically

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
